# etc/cron.d/deletedata
0 0 * * * root find /var/www/data -type f -mtime +30 -exec rm -f {} \;
